# scripts/digital-memoirs-FIXED.service
# Digital Memoirs - Event Photo Sharing Service (CORREGIDO)
# Raspberry Pi OS Bookworm (Debian 12)
#
# CAMBIOS RESPECTO A LA VERSIÓN ANTERIOR:
# - Aumentado TimeoutStartSec a 240 segundos para acomodar el delay de 180s
# - Agregado TimeoutStopSec para shutdown controlado
# - Mejorada documentación
#
# Installation:
#   sudo cp digital-memoirs-FIXED.service /etc/systemd/system/digital-memoirs.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable digital-memoirs
#   sudo systemctl start digital-memoirs
#
# View logs:
#   sudo journalctl -u digital-memoirs -f
#
# Status:
#   sudo systemctl status digital-memoirs

[Unit]
Description=Digital Memoirs - Event Photo Sharing Application
Documentation=https://github.com/yourusername/digital-memoirs
After=network-online.target dnsmasq.service
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/Downloads/repos/digital-memoirs

# SOLUCIÓN AL PROBLEMA: Aumentar timeout de systemd para acomodar el sleep de 180s
# Por defecto systemd tiene timeout de ~90s, necesitamos 240s (180s sleep + 60s margen)
TimeoutStartSec=240
TimeoutStopSec=30

# Wait 3 minutes after boot to ensure wlan0 and dnsmasq are fully ready
ExecStartPre=/bin/sleep 180

# Start Flask using uv (full path required)
# PATH is configured to include user's local bin where uv is installed
Environment="PATH=/home/pi/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/pi/.local/bin/uv run app.py

# Restart configuration
Restart=on-failure
RestartSec=10

# Logging configuration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=digital-memoirs

# Security hardening (optional - uncomment if needed)
# ProtectSystem=full
# PrivateTmp=yes
# NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
